# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for SITENAME"


start on net-device-up

stop on shutdown


respawn

setuid sluggets

chdir /home/sluggets/sites/SITENAME/source


exec ../virtualenv/bin/gunicorn --bind unix:/tmp/SITENAME.socket superlists.wsgi:application
